# lnvps_fw/lnvps_fw_service/systemd/lnvps_fw.service
[Unit]
Description=LNVPS XDP/eBPF DDoS protection firewall
Documentation=https://github.com/LNVPS/api
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/lnvps_fw_service --config /etc/lnvps_fw/config.yaml
Restart=on-failure
RestartSec=5
# Loading XDP/eBPF programs and locking map memory needs full privileges;
# run as root for broad kernel compatibility.
User=root
LimitMEMLOCK=infinity
# /var/lib/lnvps_fw holds the persisted self-signed API cert (stable fingerprint
# across restarts); /run/lnvps_fw is the root-only staging dir for self-upgrade
# .deb downloads. systemd creates both 0700, root-owned, and cleans /run on stop.
StateDirectory=lnvps_fw
StateDirectoryMode=0700
RuntimeDirectory=lnvps_fw
RuntimeDirectoryMode=0700
# Light hardening that does not interfere with the eBPF datapath.
ProtectHome=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
